# templates/researchcompendia.conf
[program:researchcompendia]
command = {{ command }}
user = {{ user }}
group = {{ group }}
autostart = true
autorestart = true
stdout_logfile = {{ logfile }}
redirect_stderr = true
killasgroup=true
